# supervisord.conf
[supervisord]
nodaemon=true

[program:jupyter-notebook]
command=jupyter notebook --ip=0.0.0.0 --port=8888 --no-browser --allow-root --NotebookApp.token='Montecristo12' /home/

[program:all]
command=/bin/bash /start-spark.sh
